# cockroach-local-cluster/haproxy/supervisord-haproxy.ini
[supervisord]
nodaemon=true

[program:haproxy]
command=/start-haproxy.sh
numprocs=1
autostart=true
autorestart=true

[program:serf]
command=/start-serf.sh
