# cockroach-local-cluster/cockroach/supervisord-cockroach.ini
[supervisord]
nodaemon=true

[program:cockroach]
command=/start-cockroach.sh
numprocs=1
autostart=true
autorestart=true

[program:serf]
command=/start-serf.sh
numprocs=1
autostart=true
autorestart=true
priority=-1
